# glogger.service
[Service]
Type=simple
UMask=0077
ExecStart=/usr/local/bin/glogger --follow

[Install]
WantedBy=multi-user.target
